# tf-kvm/compute.service
[Unit]
Description=compute-pi

[Service]
Type=simple
ExecStart=/home/ubuntu/compute/script.sh
ExecStop=/bin/kill -WINCH ${MAINPID}

